# autopraesi_python/com.autopraesi.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.autopraesi.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ehepaarschachtmaier/Desktop/Autopräsi/autopraesi_python/run_autopraesi.sh</string>
    </array>

    <!-- Jeden Donnerstag um 08:00 Uhr -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Log-Dateien -->
    <key>StandardOutPath</key>
    <string>/Users/ehepaarschachtmaier/Desktop/Autopräsi/autopraesi_python/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ehepaarschachtmaier/Desktop/Autopräsi/autopraesi_python/logs/launchd_stderr.log</string>

    <!-- Arbeitsverzeichnis -->
    <key>WorkingDirectory</key>
    <string>/Users/ehepaarschachtmaier/Desktop/Autopräsi/autopraesi_python</string>
</dict>
</plist>
